# bsp/projectroot/lib/systemd/system/mensadisplay.service
[Unit]
Description=Mensadisplay server

[Service]
ExecStart=/usr/bin/mensactrl.sh
Restart=always

[Install]
WantedBy=multi-user.target
